# systemd/oso-responder.timer
[Unit]
Description=OsO Sync — responder timer (60s poll of notes/ask/)

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Persistent=false

[Install]
WantedBy=timers.target
