# deploy/systemd-workstation/oso-memory-sync.timer
[Unit]
Description=OsO Sync memory-to-VPS timer (every 48h)

[Timer]
OnBootSec=10min
OnUnitActiveSec=48h
Unit=oso-memory-sync.service

[Install]
WantedBy=timers.target
